# deploy/signaly.service.template
[Unit]
Description=Signaly notification hub
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory={{APP_DIR}}/backend
EnvironmentFile={{APP_DIR}}/.env
ExecStart={{APP_DIR}}/.venv/bin/uvicorn main:app --host 127.0.0.1 --port 8002
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
